# themes/luci-theme-zodiacgx/root/etc/uci-defaults/30_luci-theme-zodiacgx
#!/bin/sh
uci batch <<-EOF
	set luci.themes.ZodiacGX=/luci-static/zodiacgx
	set luci.main.mediaurlbase=/luci-static/zodiacgx
	commit luci
EOF
exit 0
